# build/trace-build/lib/systemd/system/trace-api.service
[Unit]
Description=API PostgREST TRACE
After=postgresql.service
Requires=postgresql.service

[Service]
ExecStart=/usr/local/bin/postgrest /etc/trace-api.conf
Restart=always
RestartSec=3
SyslogIdentifier=trace-api

[Install]
WantedBy=multi-user.target
