# build/trace-build/DEBIAN/postinst
#!/bin/bash
set -e

echo "=== Début de l'installation de TRACE ==="

# ==============================================================================
# VARIABLES GLOBALES
# ==============================================================================
DB_NAME="parc_mobilier"
WEB_ROOT="/var/www/html/TRACE"
SCHEMA_SRC="/usr/share/trace/schema.sql"
SSL_DIR="/etc/ssl/trace"
BACKUP_DIR="/mnt/savetrace"

# Génération dynamique des secrets
DB_PASS=$(tr -dc 'A-Za-z0-9_-' < /dev/urandom | head -c 24)
JWT_SECRET=$(tr -dc 'A-Za-z0-9_-' < /dev/urandom | head -c 48)
ADMIN_EMAIL="admin@gouv.fr"
ADMIN_PASS="admin" 

# ==============================================================================
# 1. PRÉPARATION DU RÉPERTOIRE DE SAUVEGARDE
# ==============================================================================
# Création du point de montage cible. L'externalisation NFS est déléguée à 
# un paquet d'infrastructure tiers.
echo "Préparation du répertoire de sauvegarde..."
mkdir -p "$BACKUP_DIR"
# On s'assure que PostgreSQL aura les droits d'écriture une fois le dossier monté
chown postgres:postgres "$BACKUP_DIR" 2>/dev/null || true

# ==============================================================================
# 2. INITIALISATION POSTGRESQL & SCHÉMA
# ==============================================================================
echo "Initialisation de la base de données..."

pg_ctlcluster 17 main start || systemctl start postgresql || true

sudo -u postgres psql -c "SELECT pg_terminate_backend(pid) FROM pg_stat_activity WHERE datname = '$DB_NAME' AND pid <> pg_backend_pid();" > /dev/null 2>&1 || true
sudo -u postgres psql -c "DROP DATABASE IF EXISTS $DB_NAME;"
sudo -u postgres psql -c "CREATE DATABASE $DB_NAME;"

TMP_SCHEMA=$(mktemp)
cp "$SCHEMA_SRC" "$TMP_SCHEMA"
sed -i "s/__JWT_SECRET__/$JWT_SECRET/g" "$TMP_SCHEMA"
sed -i "s/__DB_PASS__/$DB_PASS/g" "$TMP_SCHEMA"
sed -i "s/__ADMIN_EMAIL__/$ADMIN_EMAIL/g" "$TMP_SCHEMA"
sed -i "s/__ADMIN_PASS__/$ADMIN_PASS/g" "$TMP_SCHEMA"

echo "Création des tables et fonctions..."
sudo -u postgres psql -d "$DB_NAME" < "$TMP_SCHEMA"
rm -f "$TMP_SCHEMA"

# ==============================================================================
# 3. IMPORT DES DONNÉES CSV & TRIGGERS
# ==============================================================================
echo "Import des référentiels CSV..."

sudo -u postgres psql -d "$DB_NAME" -c "\copy public.lieux(id, nom) FROM '/usr/share/trace/lieux.csv' DELIMITER ',' CSV HEADER;" || true
sudo -u postgres psql -d "$DB_NAME" -c "\copy public.structures(code_sages, libelle, lieu_id) FROM '/usr/share/trace/structures.csv' DELIMITER ',' CSV HEADER;" || true
sudo -u postgres psql -d "$DB_NAME" -c "SELECT setval('public.lieux_id_seq', COALESCE((SELECT MAX(id) FROM public.lieux), 1));" || true
sudo -u postgres psql -d "$DB_NAME" -c "\copy public.gabarits(reference_catalogue, categorie, nom_descriptif, caracteristiques, photo_base64) FROM '/usr/share/trace/gabarits.csv' DELIMITER ',' CSV HEADER;" || true
sudo -u postgres psql -d "$DB_NAME" -c "\copy public.mobiliers(id_metier, gabarit_id, lieu_id, code_sages, statut, remarques) FROM '/usr/share/trace/mobiliers.csv' DELIMITER ',' CSV HEADER;" || true

echo "Mise en place de la séquence d'ID Métier..."
sudo -u postgres psql -d "$DB_NAME" -c "SELECT setval('public.seq_mobilier_id', COALESCE((SELECT MAX(CAST(SUBSTRING(id_metier FROM 5) AS INTEGER)) FROM public.mobiliers), 1));" || true
sudo -u postgres psql -d "$DB_NAME" -c "CREATE TRIGGER trig_set_mob_id BEFORE INSERT ON public.mobiliers FOR EACH ROW EXECUTE FUNCTION public.set_mob_id();" || true

# ==============================================================================
# 4. CONFIGURATION DE L'API POSTGREST
# ==============================================================================
echo "Configuration du Backend PostgREST..."

cat << EOF > /etc/trace-api.conf
db-uri = "postgres://divagil:${DB_PASS}@localhost:5432/${DB_NAME}"
db-schema = "public"
db-anon-role = "divagil"
jwt-secret = "${JWT_SECRET}"
server-port = 3000
server-host = "127.0.0.1"
EOF

chmod 600 /etc/trace-api.conf

# ==============================================================================
# 5. CONFIGURATION SSL ET WEB
# ==============================================================================
echo "Génération des certificats SSL..."
mkdir -p "$SSL_DIR"
openssl req -x509 -nodes -days 365 -newkey rsa:2048 \
    -keyout "$SSL_DIR/trace.key" -out "$SSL_DIR/trace.crt" \
    -subj "/C=FR/ST=Paris/L=Paris/O=DGFIP/OU=TRACE/CN=localhost" 2>/dev/null

chown postgres:postgres "$SSL_DIR/trace.key" "$SSL_DIR/trace.crt"
chmod 600 "$SSL_DIR/trace.key"
chmod 644 "$SSL_DIR/trace.crt"
chmod 755 "$SSL_DIR"

echo "Application des droits sur le dossier Web..."
chmod -R 755 "$WEB_ROOT"
chown -R www-data:www-data "$WEB_ROOT"

ln -sf /etc/nginx/sites-available/trace /etc/nginx/sites-enabled/trace || true
rm -f /etc/nginx/sites-enabled/default || true

systemctl daemon-reload
systemctl enable trace-api
systemctl stop trace-api || true
systemctl start trace-api || true
systemctl restart nginx

echo "=== Installation de TRACE terminée avec succès ! ==="
echo "=== Identifiant : $ADMIN_EMAIL ==="
echo "=== Mot de passe : $ADMIN_PASS ==="
exit 0
